# rules.d/85-hpdrivesync.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="03f0", ATTR{idProduct}=="4607", RUN+="/home/ron/IdeaProjects/UdevRuleGen/testAction.sh"
